# apps/root-worker/nasx-root-worker.service
[Unit]
Description=NASX Root Worker
Documentation=https://github.com/nasx
After=network.target

[Service]
ExecStart=/usr/local/bin/nasx-root-worker
User=root
# Credentials are read from this file at runtime.
# Create /etc/nasx/worker.env with:
#   NATS_URL=nats://127.0.0.1:4222
#   NATS_USER=worker
#   NATS_PASS=<password matching nats.conf>
EnvironmentFile=/etc/nasx/worker.env
PrivateTmp=yes
# NoNewPrivileges must be off: the worker uses setresuid to impersonate users.
NoNewPrivileges=no
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
